# supervisord/prod_123.103.10.130.conf
[unix_http_server]
file=/tmp/hipush/supervisor.sock   

[supervisord]
logfile=/tmp/hipush/supervisord.log 
logfile_maxbytes=100MB       
logfile_backups=2            
loglevel=info                
pidfile=/tmp/hipush/supervisord.pid 
nodaemon=false               
minfds=1024                  
minprocs=200                 

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/hipush/supervisor.sock 

[group:comet]
programs=comet

[group:web]
programs=web

[group:admin]
programs=admin

[group:dog]
programs=dog

[program:comet]
command=mvn exec:java -Pcomet -Dexec.args="--runmode prod --maxconn 100000 --seq %(process_num)s --cometip 123.103.10.130 --rpcip 10.100.10.45"
process_name=%(program_name)s-%(process_num)s
numprocs=10
numprocs_start=0
autostart=true                
autorestart=false              
startsecs=3                   
redirect_stderr=true          
stdout_logfile=/data/logs/hipush/%(program_name)s-%(process_num)s.log
environment=MALLOC_ARENA_MAX='4',MAVEN_OPTS='-Xms512m -Xmx1024m',serverIdStart='0',portStart='20024'

[program:web]
command=mvn exec:java -Pweb -Dexec.args="--runmode prod --seq %(process_num)s --ip 10.100.10.45"
process_name=%(program_name)s-%(process_num)s
numprocs=2
numprocs_start=0
autostart=true                
autorestart=false             
startsecs=3                   
redirect_stderr=true          
stdout_logfile=/data/logs/hipush/%(program_name)s-%(process_num)s.log
environment=MALLOC_ARENA_MAX='4',MAVEN_OPTS='-Xms256m -Xmx512m',serverIdStart='0',portStart='21024'

[program:admin]
command=mvn exec:java -Padmin -Dexec.args="--runmode prod --seq %(process_num)s --ip 10.100.10.45"
process_name=%(program_name)s-%(process_num)s
numprocs=2
numprocs_start=0
autostart=true                
autorestart=false          
startsecs=3                   
redirect_stderr=true          
stdout_logfile=/data/logs/hipush/%(program_name)s-%(process_num)s.log
environment=MALLOC_ARENA_MAX='4',MAVEN_OPTS='-Xms256m -Xmx512m',serverIdStart='0',portStart='22024'

[program:dog]
command=mvn exec:java -Pdog -Dexec.args="--runmode prod --upstream /etc/nginx/upstream.conf"
process_name=%(program_name)s-%(process_num)s
numprocs=1
numprocs_start=0
autostart=true                
autorestart=false          
startsecs=3                   
redirect_stderr=true          
stdout_logfile=/data/logs/hipush/%(program_name)s-%(process_num)s.log
environment=MALLOC_ARENA_MAX='4',MAVEN_OPTS='-Xms256m -Xmx512m'
